# services/com.claude.usage-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.usage-collector</string>

    <key>Comment</key>
    <string>Collects API usage/cost stats for Claude Code statusline (via tkn-usage)</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{CLAUDE_HOME}}/services/usage-collector.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_HOME</key>
        <string>{{CLAUDE_HOME}}</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:{{HOME}}/.local/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{CLAUDE_HOME}}/logs/usage-collector.log</string>

    <key>StandardErrorPath</key>
    <string>{{CLAUDE_HOME}}/logs/usage-collector.error.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
