# quota-fetcher/com.claude.quota-fetcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.quota-fetcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.bun/bin/bun</string>
        <string>run</string>
        <string>{{CLAUDE_HOME}}/quota-fetcher/fetch-quota.ts</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{CLAUDE_HOME}}/quota-fetcher/fetch.log</string>
    <key>StandardErrorPath</key>
    <string>{{CLAUDE_HOME}}/quota-fetcher/fetch.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:{{HOME}}/.bun/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>
</dict>
</plist>
